# webapp.service
[Unit]
Description=Flask Web Application
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/education_app
Environment="PATH=/home/ubuntu/education_app/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PORT=5000"
ExecStart=/home/ubuntu/education_app/venv/bin/python3 /home/ubuntu/education_app/app.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
